# ansible/playbooks/deployment/templates/celery-worker.service.j2
[Unit]
Description=Celery Worker - Tenant {{ tenant_id }}
After=network.target redis-tenant-{{ tenant_id }}.service postgresql.service
Wants=redis-tenant-{{ tenant_id }}.service postgresql.service
PartOf=dotmac-tenant-{{ tenant_id }}.service

[Service]
Type=forking
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ install_dir }}

# Environment
EnvironmentFile={{ install_dir }}/.env

# Celery worker configuration
ExecStart={{ venv_dir }}/bin/celery -A dotmac.platform.celery_app worker \
    --loglevel={{ log_level | default('info') }} \
    --concurrency={{ resources.celery_workers | default('4') }} \
    --max-tasks-per-child={{ resources.celery_max_tasks_per_child | default('1000') }} \
    --time-limit={{ resources.celery_task_time_limit | default('3600') }} \
    --soft-time-limit={{ resources.celery_task_soft_time_limit | default('3000') }} \
    --pidfile={{ data_dir }}/celery-worker.pid \
    --logfile={{ data_dir }}/logs/celery-worker.log

ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
{% if resources.celery_cpu_limit is defined %}
CPUQuota={{ resources.celery_cpu_limit }}%
{% endif %}
{% if resources.celery_memory_limit is defined %}
MemoryLimit={{ resources.celery_memory_limit }}
{% endif %}
TasksMax={{ resources.celery_max_tasks | default('2048') }}

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ data_dir }} {{ install_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-worker-tenant-{{ tenant_id }}

[Install]
WantedBy=multi-user.target
